# find-master.cron
*/1 * * * * root /data/FindMaster/server_info_update.sh "/data/FindMaster"
